# scripts/linux/udiskie.service
[Unit]
Description=Udiskie automount daemon

[Install]
WantedBy=default.target

[Service]
ExecStart=/usr/bin/udiskie
