# dot_backup-tools/systemd/user/borges-backup.timer
[Unit]
Description=Run Restic Backup daily at 3 AM

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
# Small randomization to avoid exact timing collisions
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
